# sickbeard automator/lsio-arch/root/etc/s6-overlay/s6-rc.d/init-mod-sickbeard_mp4_automator-setup/up
/etc/s6-overlay/s6-rc.d/init-mod-sickbeard_mp4_automator-setup/run