# code-server-homebrew/root/etc/s6-overlay/s6-rc.d/init-mod-code-server-homebrew-add-package/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ ! "${HOMEBREW_INSTALL_COMPILERS}" = "false" ]]; then

    echo "\
        build-essential" >> /mod-repo-packages-to-install.list
fi